# crontab
# every Monday, 10 AM
0 10 * * 1 docker rmi $(docker images -f "dangling=true" -q)
